# i3/.config/systemd/user/getrxtxspeed.service
[Service]
Type=simple
ExecStart=%h/.i3/getrxtxspeed-loop 15
PIDFile=%h/.i3/.getrxtxspeed.pid
